# deploy/launchd/io.sloppad.piper-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.sloppad.piper-tts</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@VENV_DIR@@/bin/uvicorn</string>
    <string>piper_tts_server:app</string>
    <string>--app-dir</string>
    <string>@@SCRIPT_DIR@@</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8424</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PIPER_MODEL_DIR</key>
    <string>@@PIPER_MODEL_DIR@@</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/sloppad-piper-tts.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sloppad-piper-tts.log</string>
</dict>
</plist>
