# deploy/launchd/io.slopshell.codex-app-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.slopshell.codex-app-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@CODEX_PATH@@</string>
    <string>app-server</string>
    <string>--listen</string>
    <string>ws://127.0.0.1:8787</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/slopshell-codex-app-server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/slopshell-codex-app-server.log</string>
</dict>
</plist>
