# deploy/launchd/io.tabura.llm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.tabura.llm</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@LLM_SETUP_SCRIPT@@</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>TABURA_LLM_PRESET</key>
    <string>fast-qwen9b</string>
    <key>TABURA_LLM_VENV_DIR</key>
    <string>@@LLM_VENV_DIR@@</string>
    <key>TABURA_MLX_MODEL_REPO</key>
    <string>mlx-community/Qwen3.5-9B-4bit</string>
    <key>TABURA_VLLM_MLX_GIT_URL</key>
    <string>git+ssh://git@github.com/krystophny/vllm-mlx.git</string>
    <key>TABURA_VLLM_MLX_GIT_REF</key>
    <string>94fd8a68d3038eb42e5bb9d9ae9001e23ea83031</string>
    <key>TABURA_VLLM_MLX_ENABLE_BATCHING</key>
    <string>1</string>
    <key>TABURA_VLLM_MLX_USE_PAGED_CACHE</key>
    <string>1</string>
    <key>TABURA_VLLM_MLX_CACHE_MEMORY_PERCENT</key>
    <string>0.20</string>
    <key>TABURA_VLLM_MLX_TOOL_CALL_PARSER</key>
    <string>qwen</string>
    <key>TABURA_VLLM_MLX_REASONING_PARSER</key>
    <string>qwen3</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tabura-llm.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tabura-llm.log</string>
</dict>
</plist>
