# deploy/systemd/user/slopshell-codex-app-server.service
[Unit]
Description=Codex App Server (Slopshell)
After=network.target

[Service]
Type=simple
WorkingDirectory=@@REPO_ROOT@@
Environment=HOME=%h
Environment=SHELL=/usr/bin/bash
ExecStart=@@CODEX_PATH@@ app-server --listen ws://127.0.0.1:8787
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
