# deploy/systemd/user/slopshell-web.service
[Unit]
Description=Slopshell Web UI
After=network.target slopshell-codex-app-server.service slopshell-piper-tts.service helpy-mcp.service
Wants=slopshell-codex-app-server.service slopshell-piper-tts.service helpy-mcp.service

[Service]
Type=simple
WorkingDirectory=@@REPO_ROOT@@
PassEnvironment=DISPLAY WAYLAND_DISPLAY XAUTHORITY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS
Environment=HOME=%h
Environment=SHELL=/usr/bin/bash
Environment=SLOPSHELL_APP_SERVER_MODEL=gpt-5.3-codex-spark
Environment=SLOPSHELL_HELPY_SOCKET=%t/sloppy/helpy.sock
Environment=SLOPSHELL_BRAIN_GTD_SYNC=on
ExecStart=/usr/bin/bash -lc 'set -a; source ~/.bashrc >/dev/null 2>&1 || true; source ~/.config/helpy/env >/dev/null 2>&1 || true; source ~/.config/slopshell/llm.env >/dev/null 2>&1 || true; set +a; exec @@BIN_PATH@@ server --workspace-dir @@REPO_ROOT@@ --data-dir %h/.local/share/slopshell-web --control-socket %t/sloppy/control.sock --web-host @@SLOPSHELL_WEB_HOST@@ --web-port 8420 --app-server-url ws://127.0.0.1:8787 --tts-url http://127.0.0.1:8424 --dev-runtime'
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

[Install]
WantedBy=default.target
