# deploy/systemd/user/sloppad-f5-tts.service
[Unit]
Description=Sloppad F5-TTS (fast, EN+DE, GPU)
After=network.target

[Service]
Type=simple
WorkingDirectory=@@REPO_ROOT@@
Environment=HOME=%h
Environment=F5_TTS_REF_WAV=%h/.local/share/sloppad-tts/reference.wav
ExecStart=%h/.local/share/sloppad-f5-tts/venv/bin/uvicorn f5_tts_server:app --app-dir @@REPO_ROOT@@/scripts --host 127.0.0.1 --port 8424
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=default.target
