# deploy/systemd/user/tabura-stt.service
[Unit]
Description=Voxtype daemon with STT service for Tabura
After=graphical-session.target pipewire.service pipewire-pulse.service

[Service]
Type=simple
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TABURA_STT_LANGUAGE=de,en
Environment=TABURA_STT_MODEL=large-v3-turbo
ExecStart=@@REPO_ROOT@@/scripts/setup-voxtype-stt.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
